# launchd/worldcup-2026-agent.plist.example
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.local.worldcup-2026-agent</string>
  <key>ProgramArguments</key>
  <array>
    <string>/e/work/worldcup-2026-agent/run.sh</string>
  </array>
  <key>WorkingDirectory</key>
  <string>/e/work/worldcup-2026-agent</string>
  <key>StartCalendarInterval</key>
  <array>
    <dict><key>Hour</key><integer>9</integer><key>Minute</key><integer>0</integer></dict>
    <dict><key>Hour</key><integer>18</integer><key>Minute</key><integer>0</integer></dict>
  </array>
  <key>StandardOutPath</key>
  <string>/e/work/worldcup-2026-agent/logs/agent.log</string>
  <key>StandardErrorPath</key>
  <string>/e/work/worldcup-2026-agent/logs/agent.err.log</string>
</dict>
</plist>
